# supervisord.conf
[supervisord]
nodaemon=true

[program:postgres]
user=postgres
directory=/tmp
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf

[program:redis]
directory=/tmp
command=/usr/bin/redis-server

[program:resque]
directory=/root/cartodb
environment=QUEUE="imports,users,geocodings,synchronizations,tracker",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/rubies/ruby-2.2.3/bin"
command=/opt/rubies/ruby-2.2.3/bin/bundle exec /opt/rubies/ruby-2.2.3/bin/rake environment resque:work

[program:cartodb]
environment=RAILS_ENV="development",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/rubies/ruby-2.2.3/bin"
directory=/root/cartodb
command=/opt/rubies/ruby-2.2.3/bin/bundle exec rails server

[program:windshaft]
environment=RAILS_ENV="development",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/rubies/ruby-2.2.3/bin"
directory=/root/Windshaft-cartodb
command=node app.js development

[program:sqlapi]
environment=RAILS_ENV="development",PATH="/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/rubies/ruby-2.2.3/bin"
directory=/root/CartoDB-SQL-API
command=node app.js development

